# systemd/glideclaw.service
[Unit]
Description=GlideClaw personal gateway agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=glideclaw
Group=glideclaw
WorkingDirectory=/var/lib/glideclaw
EnvironmentFile=-/etc/glideclaw/glideclaw.env
ExecStart=/usr/local/bin/glideclaw run
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/glideclaw
UMask=0077
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
